# config/files/etc/systemd/system/set-user-shell-on-update.path
[Unit]
Description=Watch for ostree deployments to trigger shell update

[Path]
PathModified=/ostree/deploy/fedora/deploy
Unit=set-user-shell.service

[Install]
WantedBy=multi-user.target
